# src/async_crud_mcp/daemon/macos/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.async-crud-mcp.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>BOOTSTRAP_PATH_PLACEHOLDER</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>WORKING_DIR_PLACEHOLDER</string>

    <key>StandardOutPath</key>
    <string>STDOUT_LOG_PLACEHOLDER</string>

    <key>StandardErrorPath</key>
    <string>STDERR_LOG_PLACEHOLDER</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
